# eazybe-api/supervisord.conf
[supervisord]
nodaemon=true

[program:karafka]
command=bundle exec karafka server

[program:media]
command=rails server -b 0.0.0.0